# configs/supervisord/cleartext.carsales.com.au.conf
[program:xmpp-bot]
command = /home/admin/opt/ximpp-server/env/bin/bot /home/admin/opt/ximpp-server/configs/bot/cleartext.carsales.com.au.cfg
