# configs/supervisord/freightinvestor.cleartext.net.conf
[program:xmpp-bot]
command = /home/admin/opt/ximpp-server/env/bin/bot /home/admin/opt/ximpp-server/configs/bot/freightinvestor.cleartext.net.cfg

