# configs/supervisord/mblog.cleartext.im.conf
[program:xmpp-bot]
command = /home/admin/opt/ximpp-server/env/bin/bot /home/admin/opt/ximpp-server/configs/bot/mblog.cleartext.im.cfg

